# etc/cron.daily/oam
#!/bin/bash

if [ -z "$LANG" ] ; then
    export LANG="$(eselect locale show | tail -1| xargs)"
fi

exec /usr/bin/oam expire >>/var/log/oam/expire.log 2>&1
