# etc/gentoo-oam.cron
#!/bin/bash

exec /usr/sbin/oam expire >>/var/log/oam/expire.log
