# etc/gentoo-oam.cron
#!/bin/bash

source /etc/gentoo-oam.conf

exec /usr/sbin/oam expire >>/var/log/oam/expire.log
