# etc/oam-depclean-check.cron
#!/bin/bash

exec /usr/bin/oam-unused periodic >>/var/log/oam/expire.log 2>&1
